# site/roles/trinity/slurm/templates/systemd/minute-timer.target
[Unit]
Description=Minute Timer Target
StopWhenUnneeded=yes
