# site/roles/trinity/slurm/templates/systemd/slurm-stats.service.j2
[Unit]
Description=Collects SLURM statistics
Wants=slurm-stats.timer

[Service]
ExecStart=/usr/sbin/lslurm

[Install]
WantedBy=minute-timer.target
