# site/roles/trinity/slurm/templates/systemd/slurm-stats.timer.j2
[Unit]
Description=SLURM statistics

[Timer]
OnBootSec=5min
OnCalendar=*:0/1
Unit=minute-timer.target

[Install]
WantedBy=basic.target
