# site/roles/trinity/slurm/templates/systemd/slurm-stats.service
[Unit]
Description=Collects SLURM statistics
Wants=slurm-stats.timer

[Service]
ExecStart=/usr/local/bin/slurm.py

[Install]
WantedBy=minute-timer.target
