# site/roles/trinity/slurm/templates/systemd/slurmdbd.service.d/trinity.conf.j2
{{ ansible_managed | comment }}

[Unit]
Requires=munge.service

[Service]
Restart=always
