# site/roles/trinity/slurm/templates/systemd/slurmctld.service.d/trinity.conf.j2
{{ ansible_managed | comment }}

[Unit]
Requires=munge.service slurmdbd.service

[Service]
Restart=always
