# site/roles/trinity/slurm/templates/systemd/slurmd.service.d/trinity.conf.j2
{{ ansible_managed | comment }}

[Unit]
After=multi-user.target sssd.service
