# configuration/controller/ansible/site/roles/slurm/templates/systemd/munge.service.d/trinity.conf.j2
{{ ansible_managed | comment }}

[Unit]
After=remote-fs.target
Requires=remote-fs.target

[Service]
ExecStart=
ExecStart=/usr/sbin/munged  --key-file {{ munge_conf_path }}/munge.key

