# site/roles/trinity/kubernetes/templates/ks3-server-nodes-cleaner.timer.j2
[Unit]
Description=Run housekeeping script to delete k3s nodes that have been replaced

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target