# site/roles/trinity/aria2c/templates/aria2c.service.j2
[Unit]
Description=Aria2c systemd script

[Service]
Type=simple
ExecStart={{ trix_sbin }}/aria2c-daemon.sh
SendSIGHUP=yes

[Install]
WantedBy=multi-user.target
